# systemd/zoomfx.service
[Unit]
Description=ZoomFx service

[Service]
ExecStart=/root/zoom-rpi-switch/zoomfx.py

[Install]
WantedBy=multi-user.target
